# .files/backup.cron
0 1 * * */2 /home/skip/code/backup/backup >&2 > /var/log/backups/$(date +%Y-%m-%d).log
